# performance/gryadka/db/gryadka.conf
[program:redis]
command=/bin/bash /gryadka/run-redis.sh
autorestart=true
startretries=10

[program:gryadka]
command=/bin/bash /gryadka/run-gryadka.sh
autorestart=true
startretries=10