# performance/etcd/db/etcd.conf
[program:etcd]
command=/bin/bash /etcd/run-etcd.sh
autorestart=true
startretries=10

[program:api]
command=/bin/bash /etcd/run-tester.sh
autorestart=true
startretries=10